# bun-hono.service
[Unit]
Description=Bun Hono Application
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/bun-hono
Environment=NODE_ENV=production
Environment=PORT=8080
EnvironmentFile=/var/www/bun-hono/.env
ExecStart=/usr/local/bin/bun run server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bun-hono

[Install]
WantedBy=multi-user.target 